# ops/systemd/backup.timer
# backup.timer — Run the Journal Watch backup daily at 02:00 UTC.
#
# Persistent=true ensures the backup runs at next boot if the system was
# offline at the scheduled time (e.g. after a reboot or maintenance window).
#
# Install with: ops/systemd/install.sh

[Unit]
Description=Journal Watch — daily database backup
Requires=backup.service

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target
