# ops/systemd/backup.service
# backup.service — One-shot systemd service for Journal Watch database backups.
#
# Triggered by backup.timer.  Do not enable this unit directly.
# Install with: ops/systemd/install.sh

[Unit]
Description=Journal Watch — Restic database backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# Load all backup credentials and configuration from a single env file.
# The file must be readable only by root (chmod 600).
EnvironmentFile=/etc/restic/env

ExecStart=/opt/backup/backup.sh --verify

# Kill the backup after 30 minutes — prevents a hung pg_dump from blocking
# the next scheduled run.  Restic uploads are resumable so partial state
# is not left in the repository.
TimeoutStartSec=30min

# Write all output to the system journal.
# View with: journalctl -u backup.service [-f] [-n 100]
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backup

[Install]
WantedBy=multi-user.target
