# scripts/timers/backup.timer
[Unit]
Description=this timer runs the backup service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
